# bkp/rc.local
#!/bin/bash

export TZ=America/Sao_Paulo
env TZ=America/Sao_Paulo

sudo -H -u odroid /home/odroid/odroid-feeder/bot.sh > /tmp/bot.log 2>&1 &
sudo -H -u odroid /home/odroid/odroid-brot/brot.sh > /tmp/brot.log 2>&1 &
sudo -H -u odroid /home/odroid/bash-rest-server/index.sh > /tmp/rest.log 2>&1 &

if [ -f /aafirstboot ]; then /aafirstboot start ; fi

exit 0
